# iot/usr_lib/opkg/info/dnscrypt-proxy.list
/etc/config/dnscrypt-proxy
/usr/sbin/dnscrypt-proxy
/etc/init.d/dnscrypt-proxy
