# iot/usr_lib/opkg/info/gl-base-files.list
/rootfs-overlay/etc/lighttpd/conf.d/30-expire.conf
/rootfs-overlay/etc/chatscripts/csg_evdoipv6.chat
/rootfs-overlay/etc/chatscripts/3g.chat
/rootfs-overlay/etc/sysctl.d/local.conf
/rootfs-overlay/etc/gcom/callip.gcom
/rootfs-overlay/etc/hotplug.d/usb/20-usbmuxd
/rootfs-overlay/etc/vpn.user
/rootfs-overlay/etc/chatscripts/3gipv6.chat
/rootfs-overlay/etc/smsd.conf
/rootfs-overlay/etc/ssmtp/ssmtp.conf
/rootfs-overlay/etc/gcom/directip.gcom
/rootfs-overlay/etc/hotplug.d/tty/30-3g
/rootfs-overlay/etc/firewall.user
/rootfs-overlay/etc/hotplug.d/net/40-wwan
/rootfs-overlay/etc/uci-defaults/99-disable-upnpd
/rootfs-overlay/etc/Visible_profile
/rootfs-overlay/etc/gcom/directip-stop.gcom
/rootfs-overlay/etc/uci-defaults/101-upgrade_from_2_x_firmware
/rootfs-overlay/etc/lighttpd/conf.d/30-proxy.conf
/rootfs-overlay/etc/sysctl.d/99-disable-ipv6.conf
/rootfs-overlay/lib/network/wwan/05c6
/rootfs-overlay/etc/lighttpd/mime.conf
/rootfs-overlay/usr/share/zoneinfo/tzdata.lua
/rootfs-overlay/etc/ip-up.d/01-clear-count.sh
/rootfs-overlay/lib/upgrade/keep.d/smsforward
/rootfs-overlay/usr/bin/glautoupdater
/rootfs-overlay/etc/gcom/callip-stop.gcom
/rootfs-overlay/etc/gcom/getimsi.gcom
/rootfs-overlay/etc/gcom/ncm.json
/rootfs-overlay/lib/network/wwan/05c6:9003
/rootfs-overlay/etc/ssmtp/revaliases
/rootfs-overlay/usr/bin/switchaction
/rootfs-overlay/lib/netifd/proto/3g.sh
/rootfs-overlay/lib/netifd/dhcpv6.script
/rootfs-overlay/etc/samba/smb.conf.template
/rootfs-overlay/etc/lighttpd/conf.d/30-openssl.conf
/rootfs-overlay/etc/lighttpd/lighttpd.conf
/rootfs-overlay/etc/chatscripts/csg_3g.chat
/rootfs-overlay/lib/netifd/ppp6-up
/rootfs-overlay/etc/config/dnscrypt-proxy
/rootfs-overlay/etc/chatscripts/evdoipv6.chat
/rootfs-overlay/usr/share/dnscrypt-proxy/dnscrypt-resolvers.csv
/rootfs-overlay/etc/init.d/initswitch
/rootfs-overlay/lib/network/wwan/05c6:9215
/rootfs-overlay/lib/functions/gl_util.sh
/rootfs-overlay/etc/init.d/lighttpd
/rootfs-overlay/usr/bin/BTN_1
/rootfs-overlay/etc/uci-defaults/dhcp_server
/rootfs-overlay/usr/bin/button_pressed
/rootfs-overlay/lib/network/wwan/05c6:9090
/rootfs-overlay/usr/bin/make_wisp_relay
/rootfs-overlay/lib/netifd/proto/callip.sh
/rootfs-overlay/etc/config/glforward
/rootfs-overlay/etc/lighttpd/conf.d/30-cgi.conf
/rootfs-overlay/lib/netifd/proto/qmi.sh
/rootfs-overlay/etc/init.d/gl_safe_mode
/rootfs-overlay/etc/forward
/rootfs-overlay/etc/rc.button/BTN_0
/rootfs-overlay/etc/rc.button/BTN_1
/rootfs-overlay/etc/init.d/gl_fixdomain
/rootfs-overlay/usr/bin/wifionoff
/rootfs-overlay/lib/wifi/qcawifi_countrycode.txt
/rootfs-overlay/etc/hotplug.d/iface/17-vpn-user
/rootfs-overlay/usr/bin/button_released
/rootfs-overlay/etc/init.d/clear_flag
/rootfs-overlay/usr/bin/apstatus
/rootfs-overlay/etc/ppp/options.pptp
/rootfs-overlay/etc/rc.local
/rootfs-overlay/usr/bin/repeater_onoff
/rootfs-overlay/usr/bin/ovpnchecker.sh
/rootfs-overlay/etc/uci-defaults/99_fix_leds
/rootfs-overlay/etc/ssmtp/smsforward.upgrade
/rootfs-overlay/lib/netifd/proto/directip.sh
/rootfs-overlay/etc/hotplug.d/block/50-samba
/rootfs-overlay/etc/hotplug.d/net/30-tethering
/rootfs-overlay/etc/uci-defaults/99_fix_lan_multicast_to_unicast
/rootfs-overlay/etc/rc.button/reset
/rootfs-overlay/etc/lighttpd/server.pem
/rootfs-overlay/etc/init.d/done
/rootfs-overlay/etc/chatscripts/csg_3gipv6.chat
/rootfs-overlay/etc/init.d/smstools3
/rootfs-overlay/etc/init.d/gl_init
/rootfs-overlay/usr/bin/flashleds
/rootfs-overlay/etc/hotplug.d/net/50-ovpn
/rootfs-overlay/etc/init.d/generate_channel_list
/rootfs-overlay/usr/bin/factorydefault
/rootfs-overlay/usr/bin/sendsms
