# iot/usr_lib/opkg/info/gl-fw.list
/etc/config/glfw
/usr/lib/gl/libglfw.so
/usr/lib/libglfw.so
/usr/bin/glfw.sh
/etc/uci-defaults/glfw
/etc/init.d/glfw
