# iot/usr_lib/opkg/info/gl-repeater.list
/usr/bin/gl_health
/etc/init.d/gl_health
/etc/init.d/gl_ssids
