# iot/usr_lib/opkg/info/gl-siderouter.list
/usr/lib/libsiderouterapi.so
/etc/siderouter.fw
/etc/init.d/siderouter
/etc/config/siderouter
/etc/hotplug.d/iface/18-siderouter
/usr/bin/lancover
/usr/lib/gl/libsiderouterapi.so
