# iot/usr_lib/opkg/info/gl-tertf.list
/etc/config/tertf
/etc/init.d/gl_tertf
/etc/tertf/mac_vendor.db
/lib/upgrade/keep.d/tertf
/usr/bin/gltertf
