# iot/usr_lib/opkg/info/gl-upload.list
/etc/init.d/gl_bigdata
/usr/bin/gl_upload
/usr/bin/get_curriface.sh
/etc/config/glbigdata
