# iot/usr_lib/opkg/info/mudial.list
/usr/lib/libmudial.so
/etc/init.d/mudial
/usr/lib/gl/libmudial.so
/usr/share/mudial/mudial-mwan3.config
/etc/config/mudial
