# iot/usr_lib/opkg/info/rtty-openssl.list
/etc/config/rtty
/usr/sbin/rtty
/etc/init.d/rtty
