# iot/usr_lib/opkg/info/tor.list
/lib/upgrade/keep.d/tor
/etc/init.d/tor
/usr/sbin/tor
